# ops/backup/nails-backup.timer
[Unit]
Description=Run Nails verified backup daily

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=nails-backup.service

[Install]
WantedBy=timers.target
