# ops/backup/nails-backup.service
[Unit]
Description=Nails verified PostgreSQL backup
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/bash /opt/nails/backup/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/nails/backups /root/.hermes/profiles/nails/backups
