# ops/digest/nails-finalization-digest.timer
[Unit]
Description=Run Nails finalization digest daily

[Timer]
OnCalendar=*-*-* 23:30:00 Europe/Moscow
Persistent=true
AccuracySec=1m
Unit=nails-finalization-digest.service

[Install]
WantedBy=timers.target
